# services/kodi/run
#!/bin/bash
/sbin/setuser abc /usr/local/lib/kodi/kodi.bin --headless
