# root/etc/services.d/kodi/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc /usr/lib/kodi/kodi.bin --headless
